# templates/vlogs-backup-daily.service.j2
[Unit]
Description=VictoriaLogs daily backup to pscloud S3
After=network-online.target victoria-logs-single.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart={{ vlogs_backup_scripts_dir }}/vlogs-backup-daily.sh
